# scripts/cosmos.service
[Unit]
Description=Cosmos Server - Container Management Platform
Documentation=https://github.com/namastexlabs/Cosmos-Server-1
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cosmos
Group=cosmos
WorkingDirectory=/opt/cosmos
ExecStart=/opt/cosmos/bin/cosmos
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=append:/opt/cosmos/logs/cosmos.log
StandardError=append:/opt/cosmos/logs/cosmos-error.log

# Environment Variables
Environment=COSMOS_CONFIG_FOLDER=/opt/cosmos/config
Environment=COSMOS_HTTP_PORT=8080
Environment=COSMOS_HTTPS_MODE=DISABLED

# For Proxmox LXC runtime - allow API access
# Uncomment and configure these for Proxmox mode:
# Environment=COSMOS_RUNTIME_TYPE=proxmox
# Environment=COSMOS_PROXMOX_HOST=10.114.1.1:8006
# Environment=COSMOS_PROXMOX_NODE=dl380
# Environment=COSMOS_PROXMOX_TOKEN_ID=genieos1@pve!api-token
# Environment=COSMOS_PROXMOX_TOKEN_SECRET=<your-token-secret>
# Environment=COSMOS_PROXMOX_STORAGE=local-lvm
# Environment=COSMOS_PROXMOX_SKIP_TLS=true

# Security hardening (comment out if causing issues)
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/cosmos/config /opt/cosmos/data /opt/cosmos/logs

# Resource limits (adjust as needed)
# MemoryMax=2G
# CPUQuota=200%

[Install]
WantedBy=multi-user.target
